# systemd/llm2ssh-bot.service
[Unit]
Description=llm2ssh Telegram bot
Documentation=man:llm2ssh
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/llm2ssh/bot.env

[Service]
User=llm2ssh-bot
Group=llm2ssh
EnvironmentFile=/etc/llm2ssh/bot.env
ExecStart=/usr/local/lib/llm2ssh/bot/llm2ssh-botd
Restart=always
RestartSec=5
# The bot's whole job is to run sudo (status/docker/freeze/grant, and relay as
# an agent user), and grants WRITE /etc/sudoers.d. So NoNewPrivileges must be no,
# and filesystem sandboxing is intentionally minimal — aggressive ProtectSystem/
# ProtectHome would break sudo, the relay (needs /home/<agent>), and grants. The
# real boundary is the bot's tight sudoers whitelist, not systemd confinement.
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
